# lib/scheduler/deploy/ceo-schedulerd.service
# ceo-schedulerd — systemd USER unit template (Linux/WSL keep-alive, #142).
#
# This is a TEMPLATE: edit the paths and host id below, then install as a user
# unit. Smoke-tested live on ML-1 (WSL) 2026-06-28: with `bun install` run and
# WorkingDirectory set, the daemon comes up healthy and heartbeats; the
# Restart=always respawn floor is now exercised on hardware.
#
# Install (user scope, no root):
#   ( cd ~/code/claude-ceo/lib/scheduler && bun install )  # REQUIRED: creates node_modules (croner); without it the daemon crash-loops on ENOENT
#   mkdir -p ~/.config/systemd/user
#   cp ceo-schedulerd.service ~/.config/systemd/user/
#   # edit the copy: ExecStart, WorkingDirectory, CEO_VAULT, CEO_HOSTNAME
#   systemctl --user daemon-reload
#   systemctl --user enable --now ceo-schedulerd.service
#   loginctl enable-linger "$USER"   # keep it running when you're logged out
#   systemctl --user status ceo-schedulerd
#   journalctl --user -u ceo-schedulerd -f
#
# Verify liveness independently of systemd:  ceo doctor   (reads the heartbeat).

[Unit]
Description=CEO playbook scheduler daemon (ceo-schedulerd)
After=network-online.target
Wants=network-online.target
# Cap respawn storms: at most 10 restarts per 5 min, then stop and surface it.
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
# EDIT: this repo's lib/scheduler dir. REQUIRED — the daemon resolves its `@/*`
# tsconfig path alias (and node_modules) relative to the working directory; run
# from elsewhere, `bun run` cannot resolve `@/cron` and the unit crash-loops.
WorkingDirectory=%h/code/claude-ceo/lib/scheduler
# EDIT: absolute path to bun and to this repo's lib/scheduler/src/main.ts.
ExecStart=%h/.bun/bin/bun run %h/code/claude-ceo/lib/scheduler/src/main.ts
# EDIT: the synced Obsidian vault root that contains CEO/registry.json.
Environment=CEO_VAULT=%h/Documents/Obsidian
# REQUIRED: this host's unique CEO id (install.md §2). systemd does NOT source
# ~/.ceo/config, so the daemon reads CEO_HOSTNAME only from here. Left as the
# CHANGEME sentinel, the host registers literally as "CHANGEME" — it does NOT
# fall back; delete this line entirely to fall back to `hostname -s`.
Environment=CEO_HOSTNAME=CHANGEME
# EDIT (optional): pin the dispatch binary if ceo-cron.sh is not on PATH.
# Environment=CEO_CRON_BIN=%h/code/claude-ceo/scripts/ceo-cron.sh
Restart=always
RestartSec=10
# Give SIGTERM time for the loop to finish its current tick and exit cleanly.
TimeoutStopSec=70

[Install]
WantedBy=default.target
